# services/daemon/install/macos/pkg/postinstall
#!/usr/bin/env bash
# postinstall — runs after pkgbuild places files at their install locations.
#
# What this script does:
#   1. Clears the macOS quarantine extended attribute on the daemon binary so
#      Gatekeeper does not block the first run with an "unidentified developer"
#      error.  (ADR-020, Wave 0 confirmed decision 2026-05-09)
#   2. CRITICAL (ADR-022 Constraint 1): Detects and bootout/unloads the OLD
#      com.mtga-companion.daemon LaunchAgent if it is registered, BEFORE
#      registering the new one.  Prevents two daemon processes from running
#      simultaneously after an in-place upgrade.
#   3. Writes a default LaunchAgent plist to ~/Library/LaunchAgents/ using the
#      new label com.vaultmtg.daemon (per-user, not system-wide — no sudo).
#   4. Loads the new LaunchAgent with launchctl bootstrap so the daemon starts
#      immediately without requiring a logout/login cycle.
#
# The binary is installed by pkgbuild into /usr/local/bin/vaultmtg-daemon
# (the package root maps to / at install time).
#
# The Clerk env vars (CLERK_FRONTEND_API, CLERK_PUBLISHABLE_KEY) and
# cloud_api_url are baked in at build time via the CI workflow's env block
# and embedded into the plist EnvironmentVariables and daemon.json respectively.

set -euo pipefail

INSTALL_DIR="/usr/local/bin"
BINARY_NAME="vaultmtg-daemon"
BINARY_PATH="${INSTALL_DIR}/${BINARY_NAME}"

# Share directory constant — must match SHARE_DIR in build-pkg.sh.
# The uninstall script is installed here by build-pkg.sh; we echo the path to
# the user at the end of postinstall so they know how to remove the daemon.
SHARE_DIR="/usr/local/share/vaultmtg"

# ADR-022 Phase 2: new label (production).
PLIST_LABEL="com.vaultmtg.daemon"
# Legacy label — detected and unloaded before registering the new one.
PLIST_LABEL_LEGACY="com.mtga-companion.daemon"

# Resolve the real (non-root) user so all files under $HOME are created with
# correct ownership. postinstall runs as root; $SUDO_USER names the invoking
# user. Fall back to the stat owner of $HOME for non-interactive installs.
REAL_USER="${SUDO_USER:-$(stat -f '%Su' "${HOME}")}"
REAL_UID=$(id -u "${REAL_USER}")
REAL_HOME=$(eval echo "~${REAL_USER}")

PLIST_DIR="${REAL_HOME}/Library/LaunchAgents"
PLIST_PATH="${PLIST_DIR}/${PLIST_LABEL}.plist"
PLIST_PATH_LEGACY="${PLIST_DIR}/${PLIST_LABEL_LEGACY}.plist"

# ADR-022 Phase 2: new config dir (~/.vaultmtg).
CONFIG_DIR="${REAL_HOME}/.vaultmtg"
CONFIG_FILE="${CONFIG_DIR}/daemon.json"

LOG_FILE="${REAL_HOME}/Library/Logs/vaultmtg-daemon.log"

# Values are substituted at build time by .github/workflows/daemon-release.yml
# using `sed` BEFORE pkgbuild is invoked. pkgbuild does NOT propagate env vars
# from CI to the installed postinstall script — the postinstall runs in macOS
# Installer's own (empty) environment on the user's machine, so `${VAR:-default}`
# expansions here would always resolve to the default. Substituting placeholders
# at build time is the only reliable way to inject build-specific values.
CLOUD_API_URL="__VAULTMTG_CLOUD_API_URL__"
CLERK_FRONTEND_API="__CLERK_FRONTEND_API__"
CLERK_PUBLISHABLE_KEY="__CLERK_PUBLISHABLE_KEY__"
CLERK_OAUTH_CLIENT_ID="__CLERK_OAUTH_CLIENT_ID__"

# Fail loudly if any placeholder was not substituted at build time.
for var in CLOUD_API_URL CLERK_FRONTEND_API CLERK_PUBLISHABLE_KEY CLERK_OAUTH_CLIENT_ID; do
  eval "value=\${$var}"
  case "$value" in
    __*__)
      echo "[postinstall] ERROR: build-time substitution did not run for ${var} (still '${value}')" >&2
      exit 1
      ;;
  esac
done

# ---------------------------------------------------------------------------
# 1. Clear quarantine attribute (ADR-020 §macOS Quarantine Clearing).
#
#    -dr = dereference symlinks and recurse (safe on a single binary).
#    This must run before the LaunchAgent loads the binary, otherwise macOS
#    will block execution on first run with a Gatekeeper "unidentified developer"
#    dialog — even for signed binaries that have not yet been notarized.
# ---------------------------------------------------------------------------
echo "[postinstall] clearing quarantine attribute on ${BINARY_PATH}"
xattr -dr com.apple.quarantine "${BINARY_PATH}" || true

# ---------------------------------------------------------------------------
# 2. Write daemon.json with cloud_api_url (required field).
#    Only written on a fresh install — skip if the file already exists so
#    re-installs do not clobber the user's existing api_key / daemon_jwt.
#    Use install(1) so file ownership is set to the real user, not root.
# ---------------------------------------------------------------------------
install -d -m 755 -o "${REAL_USER}" "${CONFIG_DIR}"
if [[ ! -f "${CONFIG_FILE}" ]]; then
  echo "[postinstall] writing ${CONFIG_FILE}"
  install -m 600 -o "${REAL_USER}" /dev/null "${CONFIG_FILE}"
  cat > "${CONFIG_FILE}" <<JSON
{
  "cloud_api_url": "${CLOUD_API_URL}",
  "sync_enabled": true
}
JSON
fi

# ---------------------------------------------------------------------------
# 3. Ensure LaunchAgents directory exists (owned by the real user).
# ---------------------------------------------------------------------------
install -d -m 755 -o "${REAL_USER}" "${PLIST_DIR}"

# ---------------------------------------------------------------------------
# 4. CRITICAL (ADR-022 Constraint 1): Stop and unregister the OLD daemon label
#    BEFORE writing and bootstrapping the new plist.
#
#    If the old label (com.mtga-companion.daemon) is still registered with
#    launchd, both daemons would run simultaneously after this postinstall
#    completes — causing duplicate log ingestion and BFF event duplication.
#
#    Steps:
#      a. Try label-based bootout (modern, macOS 10.11+) as the real user.
#      b. Fallback: file-based bootout using the legacy plist path.
#      c. Remove the legacy plist so it is not re-loaded at next login.
#    All steps use `|| true` — a fresh install has no legacy label registered.
# ---------------------------------------------------------------------------
echo "[postinstall] checking for legacy launchd label ${PLIST_LABEL_LEGACY}"
if launchctl asuser "${REAL_UID}" launchctl list "${PLIST_LABEL_LEGACY}" >/dev/null 2>&1; then
  echo "[postinstall] found legacy daemon — booting out ${PLIST_LABEL_LEGACY}"
  launchctl asuser "${REAL_UID}" launchctl bootout "gui/${REAL_UID}/${PLIST_LABEL_LEGACY}" 2>/dev/null || \
    launchctl asuser "${REAL_UID}" launchctl bootout "gui/${REAL_UID}" "${PLIST_PATH_LEGACY}" 2>/dev/null || true
  echo "[postinstall] legacy daemon stopped"
else
  echo "[postinstall] no legacy daemon found — skipping (fresh install)"
fi

# Remove legacy plist (idempotent — ignore missing-file errors).
if [[ -f "${PLIST_PATH_LEGACY}" ]]; then
  echo "[postinstall] removing legacy plist: ${PLIST_PATH_LEGACY}"
  rm -f "${PLIST_PATH_LEGACY}"
fi

# ---------------------------------------------------------------------------
# 5. Write the new LaunchAgent plist (com.vaultmtg.daemon).
#    - RunAtLoad=true   start at login.
#    - KeepAlive=true   restart if the daemon exits unexpectedly.
#    - EnvironmentVariables injects the Clerk values the PKCE flow requires.
#
#    Use install(1) to pre-create the file owned by the real user before
#    writing content. LaunchAgent plists must be owned by the user, not root;
#    a root-owned plist causes launchd to refuse to load it.
# ---------------------------------------------------------------------------
echo "[postinstall] writing LaunchAgent plist to ${PLIST_PATH}"
install -m 644 -o "${REAL_USER}" /dev/null "${PLIST_PATH}"
cat > "${PLIST_PATH}" <<PLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>${PLIST_LABEL}</string>

    <key>ProgramArguments</key>
    <array>
        <string>${BINARY_PATH}</string>
        <string>-config</string>
        <string>${CONFIG_FILE}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLERK_FRONTEND_API</key>
        <string>${CLERK_FRONTEND_API}</string>
        <key>CLERK_PUBLISHABLE_KEY</key>
        <string>${CLERK_PUBLISHABLE_KEY}</string>
        <key>CLERK_OAUTH_CLIENT_ID</key>
        <string>${CLERK_OAUTH_CLIENT_ID}</string>
        <key>VAULTMTG_DAEMON_CLOUD_API_URL</key>
        <string>${CLOUD_API_URL}</string>
    </dict>

    <!-- Start automatically when the user logs in. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart the daemon if it exits for any reason. -->
    <key>KeepAlive</key>
    <true/>

    <!-- Minimum seconds between restarts to prevent launchd restart storms. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>${LOG_FILE}</string>
    <key>StandardErrorPath</key>
    <string>${LOG_FILE}</string>
</dict>
</plist>
PLIST

# ---------------------------------------------------------------------------
# 6. Load the new LaunchAgent into the correct user gui domain.
#
#    pkgbuild postinstall scripts run as root. LaunchAgents must be loaded
#    into the user's gui domain (not the system domain). Use launchctl
#    bootstrap/bootout with the gui/<uid> target, run as the actual user
#    via launchctl asuser so the agent is registered correctly.
# ---------------------------------------------------------------------------
echo "[postinstall] loading LaunchAgent into gui/${REAL_UID} (user: ${REAL_USER})"

# Re-enable the label before bootstrapping. If a previous uninstall called
# "launchctl bootout", launchd marks the label as disabled and a subsequent
# "launchctl bootstrap" fails with error 5 (I/O error). Enabling first
# clears that disabled flag regardless of whether it was set.
launchctl enable "gui/${REAL_UID}/${PLIST_LABEL}"

# Stop any still-running daemon process (covers the same-label reinstall case).
launchctl asuser "${REAL_UID}" launchctl bootout "gui/${REAL_UID}/${PLIST_LABEL}" 2>/dev/null || \
    launchctl asuser "${REAL_UID}" launchctl bootout "gui/${REAL_UID}" "${PLIST_PATH}" 2>/dev/null || true

# Kill any stale processes that slipped past launchctl (e.g. orphaned after a
# crash or a prior install that used a different label).
pkill -x vaultmtg-daemon 2>/dev/null || true

# Brief wait so the process has time to exit cleanly before bootstrap starts.
sleep 1

# Wrap bootstrap in a conditional so a transient launchd failure (e.g. when
# the GUI session has not yet started in a script-driven install) does NOT
# cause the entire pkg install to roll back. RunAtLoad=true guarantees the
# daemon will start at the user's next login even if bootstrap is skipped.
if launchctl asuser "${REAL_UID}" launchctl bootstrap "gui/${REAL_UID}" "${PLIST_PATH}"; then
  echo "[postinstall] daemon bootstrapped and running"
else
  echo "[postinstall] WARNING: launchctl bootstrap returned non-zero (GUI session may not be active)."
  echo "[postinstall] The daemon is registered with RunAtLoad=true and will start automatically at next login."
fi

echo "[postinstall] done"
echo "  Binary : ${BINARY_PATH}"
echo "  plist  : ${PLIST_PATH}"
echo "  Logs   : ${LOG_FILE}"
echo "To uninstall: sudo ${SHARE_DIR}/uninstall.sh"
