# services/daemon/install/macos/pkg/postinstall
#!/usr/bin/env bash
# postinstall — runs after pkgbuild places files at their install locations.
#
# What this script does:
#   1. Clears the macOS quarantine extended attribute on the daemon binary so
#      Gatekeeper does not block the first run with an "unidentified developer"
#      error.  (ADR-020, Wave 0 confirmed decision 2026-05-09)
#   2. Writes a default LaunchAgent plist to ~/Library/LaunchAgents/ (per-user,
#      not system-wide — no LaunchDaemon, no sudo for the agent itself).
#   3. Loads the LaunchAgent with launchctl so the daemon starts immediately
#      without requiring a logout/login cycle.
#
# The binary is installed by pkgbuild into /usr/local/bin/vaultmtg-daemon
# (the package root maps to / at install time).
#
# The Clerk env vars (CLERK_FRONTEND_API, CLERK_PUBLISHABLE_KEY) and
# cloud_api_url are baked in at build time via the CI workflow's env block
# and embedded into the plist EnvironmentVariables and daemon.json respectively.

set -euo pipefail

INSTALL_DIR="/usr/local/bin"
BINARY_NAME="vaultmtg-daemon"
BINARY_PATH="${INSTALL_DIR}/${BINARY_NAME}"

PLIST_LABEL="com.mtga-companion.daemon"

# Resolve the real (non-root) user so all files under $HOME are created with
# correct ownership. postinstall runs as root; $SUDO_USER names the invoking
# user. Fall back to the stat owner of $HOME for non-interactive installs.
REAL_USER="${SUDO_USER:-$(stat -f '%Su' "${HOME}")}"
REAL_UID=$(id -u "${REAL_USER}")
REAL_HOME=$(eval echo "~${REAL_USER}")

PLIST_DIR="${REAL_HOME}/Library/LaunchAgents"
PLIST_PATH="${PLIST_DIR}/${PLIST_LABEL}.plist"

CONFIG_DIR="${REAL_HOME}/.mtga-companion"
CONFIG_FILE="${CONFIG_DIR}/daemon.json"

LOG_FILE="${REAL_HOME}/Library/Logs/vaultmtg-daemon.log"

# Values are substituted at build time by .github/workflows/daemon-release.yml
# using `sed` BEFORE pkgbuild is invoked. pkgbuild does NOT propagate env vars
# from CI to the installed postinstall script — the postinstall runs in macOS
# Installer's own (empty) environment on the user's machine, so `${VAR:-default}`
# expansions here would always resolve to the default. Substituting placeholders
# at build time is the only reliable way to inject build-specific values.
CLOUD_API_URL="__VAULTMTG_CLOUD_API_URL__"
CLERK_FRONTEND_API="__CLERK_FRONTEND_API__"
CLERK_PUBLISHABLE_KEY="__CLERK_PUBLISHABLE_KEY__"
CLERK_OAUTH_CLIENT_ID="__CLERK_OAUTH_CLIENT_ID__"

# Fail loudly if any placeholder was not substituted at build time.
for var in CLOUD_API_URL CLERK_FRONTEND_API CLERK_PUBLISHABLE_KEY CLERK_OAUTH_CLIENT_ID; do
  eval "value=\${$var}"
  case "$value" in
    __*__)
      echo "[postinstall] ERROR: build-time substitution did not run for ${var} (still '${value}')" >&2
      exit 1
      ;;
  esac
done

# ---------------------------------------------------------------------------
# 1. Clear quarantine attribute (ADR-020 §macOS Quarantine Clearing).
#
#    -dr = dereference symlinks and recurse (safe on a single binary).
#    This must run before the LaunchAgent loads the binary, otherwise macOS
#    will block execution on first run with a Gatekeeper "unidentified developer"
#    dialog — even for signed binaries that have not yet been notarized.
# ---------------------------------------------------------------------------
echo "[postinstall] clearing quarantine attribute on ${BINARY_PATH}"
xattr -dr com.apple.quarantine "${BINARY_PATH}" || true

# ---------------------------------------------------------------------------
# 2. Write daemon.json with cloud_api_url (required field).
#    Only written on a fresh install — skip if the file already exists so
#    re-installs do not clobber the user's existing api_key / daemon_jwt.
#    Use install(1) so file ownership is set to the real user, not root.
# ---------------------------------------------------------------------------
install -d -m 755 -o "${REAL_USER}" "${CONFIG_DIR}"
if [[ ! -f "${CONFIG_FILE}" ]]; then
  echo "[postinstall] writing ${CONFIG_FILE}"
  install -m 600 -o "${REAL_USER}" /dev/null "${CONFIG_FILE}"
  cat > "${CONFIG_FILE}" <<JSON
{
  "cloud_api_url": "${CLOUD_API_URL}",
  "sync_enabled": true
}
JSON
fi

# ---------------------------------------------------------------------------
# 3. Ensure LaunchAgents directory exists (owned by the real user).
# ---------------------------------------------------------------------------
install -d -m 755 -o "${REAL_USER}" "${PLIST_DIR}"

# ---------------------------------------------------------------------------
# 4. Write the LaunchAgent plist.
#    - RunAtLoad=true   start at login.
#    - KeepAlive=true   restart if the daemon exits unexpectedly.
#    - EnvironmentVariables injects the Clerk values the PKCE flow requires.
#
#    Use install(1) to pre-create the file owned by the real user before
#    writing content. LaunchAgent plists must be owned by the user, not root;
#    a root-owned plist causes launchd to refuse to load it.
# ---------------------------------------------------------------------------
echo "[postinstall] writing LaunchAgent plist to ${PLIST_PATH}"
install -m 644 -o "${REAL_USER}" /dev/null "${PLIST_PATH}"
cat > "${PLIST_PATH}" <<PLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>${PLIST_LABEL}</string>

    <key>ProgramArguments</key>
    <array>
        <string>${BINARY_PATH}</string>
        <string>-config</string>
        <string>${CONFIG_FILE}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLERK_FRONTEND_API</key>
        <string>${CLERK_FRONTEND_API}</string>
        <key>CLERK_PUBLISHABLE_KEY</key>
        <string>${CLERK_PUBLISHABLE_KEY}</string>
        <key>CLERK_OAUTH_CLIENT_ID</key>
        <string>${CLERK_OAUTH_CLIENT_ID}</string>
        <key>MTGA_DAEMON_CLOUD_API_URL</key>
        <string>${CLOUD_API_URL}</string>
    </dict>

    <!-- Start automatically when the user logs in. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart the daemon if it exits for any reason. -->
    <key>KeepAlive</key>
    <true/>

    <!-- Minimum seconds between restarts to prevent launchd restart storms. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>${LOG_FILE}</string>
    <key>StandardErrorPath</key>
    <string>${LOG_FILE}</string>
</dict>
</plist>
PLIST

# ---------------------------------------------------------------------------
# 5. Load the LaunchAgent into the correct user gui domain.
#
#    pkgbuild postinstall scripts run as root. LaunchAgents must be loaded
#    into the user's gui domain (not the system domain). Use launchctl
#    bootstrap/bootout with the gui/<uid> target, run as the actual user
#    via launchctl asuser so the agent is registered correctly.
#
#    $SUDO_USER is set by sudo to the invoking user's name; if absent
#    (direct root shell) fall back to the owner of HOME.
# ---------------------------------------------------------------------------
REAL_USER="${SUDO_USER:-$(stat -f '%Su' "${HOME}")}"
REAL_UID=$(id -u "${REAL_USER}")

echo "[postinstall] loading LaunchAgent into gui/${REAL_UID} (user: ${REAL_USER})"
# Bootout any previous registration so bootstrap starts from a clean slate.
launchctl asuser "${REAL_UID}" launchctl bootout "gui/${REAL_UID}" "${PLIST_PATH}" 2>/dev/null || true

# Re-enable the label before bootstrapping. If a previous uninstall called
# "launchctl bootout", launchd marks the label as disabled and a subsequent
# "launchctl bootstrap" fails with error 5 (I/O error). Enabling first
# clears that disabled flag regardless of whether it was set.
launchctl enable "gui/${REAL_UID}/${PLIST_LABEL}"

# Wrap bootstrap in a conditional so a transient launchd failure (e.g. when
# the GUI session has not yet started in a script-driven install) does NOT
# cause the entire pkg install to roll back. RunAtLoad=true guarantees the
# daemon will start at the user's next login even if bootstrap is skipped.
if launchctl asuser "${REAL_UID}" launchctl bootstrap "gui/${REAL_UID}" "${PLIST_PATH}"; then
  echo "[postinstall] daemon bootstrapped and running"
else
  echo "[postinstall] WARNING: launchctl bootstrap returned non-zero (GUI session may not be active)."
  echo "[postinstall] The daemon is registered with RunAtLoad=true and will start automatically at next login."
fi

echo "[postinstall] done"
echo "  Binary : ${BINARY_PATH}"
echo "  plist  : ${PLIST_PATH}"
echo "  Logs   : ${LOG_FILE}"
