# .github/deploy/mtga-sync.service
[Unit]
Description=MTGA Sync Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=mtga
ExecStartPre=/bin/bash -c 'export DATABASE_URL=$(aws ssm get-parameter --name /mtga/prod/database_url --with-decryption --query Parameter.Value --output text --region us-east-1) && echo "DATABASE_URL=${DATABASE_URL}" > /etc/mtga-sync.env'
EnvironmentFile=/etc/mtga-sync.env
Environment=SYNC_REFRESH_HOUR=2
ExecStart=/usr/local/bin/mtga-sync
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mtga-sync

[Install]
WantedBy=multi-user.target
